# projects/jfnetboot/jfnetboot.service
[Unit]
Description=JavaForce Net Boot Service

[Service]
ExecStart=/usr/bin/jfnetboot
Restart=always

[Install]
