# projects/jfsocks/jfsockssvc.service
[Unit]
Description=JavaForce Socks Server
After=jfbusserver.service

[Service]
ExecStart=/usr/bin/jfsocks
Restart=always

[Install]
