# projects/jfbusserver/jfbusserver.service
[Unit]
Description=JavaForce Bus Server
After=

[Service]
ExecStart=/usr/sbin/jfbusserver
Restart=always

[Install]
WantedBy=multi-user.target
