# projects/jfdhcp/jfdhcp.service
[Unit]
Description=JavaForce DHCP Server
After=jfbusserver.service

[Service]
ExecStart=/usr/bin/jfdhcpsvc
Restart=always

[Install]
