# projects/jfqemu/jfqemu.service
[Unit]
Description=JavaForce QEMU Server
After=jfbusserver.service

[Service]
ExecStart=/usr/bin/jfqemu
Restart=always

[Install]
