# projects/jfcontrols/jfcontrols.service
[Unit]
Description=jfControls

[Service]
ExecStart=/usr/bin/jfcontrols
Restart=always
