# projects/jfssh/jfssh.service
[Unit]
Description=JavaForce SSH Server
After=jfbusserver.service

[Service]
ExecStart=/usr/bin/jfsshsvc
ExecStop=/usr/bin/jfsshsvc --stop
Restart=always

[Install]
WantedBy=multi-user.target
