# projects/jpbxlite/jpbxlite.service
[Unit]
Description=JavaForce PBX Server
After=jfbusserver.service

[Service]
ExecStart=/usr/bin/jpbxlite
Restart=always

[Install]
