# examples/service/service.service
[Unit]
Description=JavaForce Example Service
After=jfbusserver.service

[Service]
ExecStart=/usr/bin/service
ExecStop=/usr/bin/service --stop
ExecReload=/usr/bin/service --reload
Restart=always

[Install]
WantedBy=multi-user.target
