# projects/jfnetworkmgr/jfnetworkmgr.service
[Unit]
Description=JavaForce Network Manager
After=jfbusserver.service

[Service]
ExecStart=/usr/bin/jfnetworkmgr
Restart=always

[Install]
