# projects/jflogon/deb/postinst
#!/bin/bash

set -e

systemctl enable jflogon

